# debian/pop-upgrade.postinst
#!/bin/sh
set -x
systemctl --global enable pop-upgrade-notify.timer
systemctl disable pop-upgrade

if ! test -e /var/lib/pop-upgrade/restarting && ! test -e /pop-upgrade && ! test -e /pop_preparing_release_upgrade; then
    if systemctl is-active pop-upgrade; then
        systemctl restart pop-upgrade
    fi
fi

#DEBHELPER#

exit 0
