# debian/pop-upgrade.postinst
#!/bin/sh
set -x
systemctl --global enable pop-upgrade-notify.timer

#DEBHELPER#

exit 0
